# schedules/dotfiles-sync.timer
[Unit]
Description=Dotfiles periodic sync timer

[Timer]
OnCalendar=daily
Persistent=true
Unit=dotfiles-sync.service

[Install]
WantedBy=timers.target
